# cfg/pi-gen/stage5/02-drt-wifi/files/add-drt-wifi.service
[Unit]
Description=Add DRT wifi networks
Before=networking.service

[Service]
Type=oneshot
ExecStart=/opt/drt-wifi/add-drt-wifi.sh

[Install]
RequiredBy=network.target

